# files/files/etc/profile.d/charlisim-base.sh
export EDITOR='vim'

alias vi=/usr/bin/vim
